# launchd/com.engram.v3.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
        "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Engram V3 watchdog.

  Install:
    cp launchd/com.engram.v3.watchdog.plist ~/Library/LaunchAgents/
    # Edit ProgramArguments[0] to the absolute path for scripts/engram_watchdog.sh.
    launchctl load ~/Library/LaunchAgents/com.engram.v3.watchdog.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.engram.v3.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/REPLACE/WITH/ABSOLUTE/PATH/TO/engram/scripts/engram_watchdog.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/engram.watchdog.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/engram.watchdog.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
