# launchd/com.engram.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
        "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Engram user-launchd manifest.

  Install:
    engram setup
    # or: scripts/install_launchd.sh

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.engram.bridge.plist
    rm ~/Library/LaunchAgents/com.engram.bridge.plist

  Logs:
    tail -f /tmp/engram.bridge.out.log
    tail -f /tmp/engram.bridge.err.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.engram.bridge</string>

    <!-- EDIT: absolute path to `uv`.
         Get it with: `which uv` after install.sh finishes. -->
    <key>ProgramArguments</key>
    <array>
        <string>/REPLACE/WITH/ABSOLUTE/PATH/TO/uv</string>
        <string>run</string>
        <string>--project</string>
        <string>/REPLACE/WITH/ABSOLUTE/PATH/TO/engram-repo</string>
        <string>python</string>
        <string>-m</string>
        <string>engram.main</string>
    </array>

    <!-- EDIT: absolute path to the repo checkout. -->
    <key>WorkingDirectory</key>
    <string>/REPLACE/WITH/ABSOLUTE/PATH/TO/engram-repo</string>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Keep in sync with the deployed bridge env. HOME is explicit so
             ~/.engram resolves the same way under launchd and a login shell. -->
        <key>PATH</key>
        <string>/REPLACE/WITH/OPTIONAL/NODE/PATH/PREFIX//REPLACE/WITH/HOME/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>HOME</key>
        <string>/REPLACE/WITH/HOME</string>
        <key>ENGRAM_ENV_FILE</key>
        <string>/REPLACE/WITH/ABSOLUTE/PATH/TO/engram.env</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <!-- Avoid tight crash loops: wait 30s before restart. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>

    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>8192</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/engram.bridge.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/engram.bridge.err.log</string>

    <!-- Be a polite user agent: low I/O priority, normal CPU. -->
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
